# services/lowbat-suspend.service
[Unit]
Description=Suspend on low battery

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/lowbat-suspend.sh
