# services/lowbat-suspend.service
[Unit]
Description=Suspend on low battery

[Service]
Type=oneshot
ExecStart=/home/esteban/dotfiles/services/lowbat-suspend.sh
User=root
