# services/lowbat-suspend.timer
[Unit]
Description=1 minute timer for lowbat service

[Timer]
OnBootSec=5min
OnCalendar=*:*:0

[Install]
WantedBy=timer.target
