# services/.config/systemd/user/lowbat-suspend.service
[Unit]
Description=Suspend on low battery

[Service]
Type=oneshot
ExecStart=./lowbat-suspend.sh
User=root
